# utils/qmi-modem-410-init/files/qmi-modem-410-init.init
#!/bin/sh /etc/rc.common
# Init script for QMI Modem Initialization

START=80
STOP=10

USE_PROCD=1
PROG=/usr/bin/qmi-modem-init

start_service() {
    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param stdout 1
    procd_set_param stderr 1
    # Delay execution to ensure network is ready
    # procd_set_param start_delay 30
    procd_close_instance
}

stop_service() {
    # Nothing to stop
    return 0
}
